# launchd/com.hrp.quality-sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hrp.quality-sentinel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fer/Projects/HRP/.venv/bin/python</string>
        <string>-m</string>
        <string>hrp.agents.run_job</string>
        <string>--job</string>
        <string>quality-sentinel</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/fer/Projects/HRP</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/fer/hrp-data/logs/quality-sentinel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fer/hrp-data/logs/quality-sentinel.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/fer/Projects/HRP</string>
    </dict>
</dict>
</plist>
